# deploy/usa-wa-succession-invariants.service
[Unit]
Description=usa-wa — succession invariants (chamber counts + seat occupancy, #107) one-shot
# After the three daily refreshes so it checks the freshly-rebuilt open-seat cohort — a
# ghost-open predecessor (missing operator `departed`/`vacated`) or an unfilled seat (missing
# `seated`) shows up the morning after the wire re-drive. Ordering only — a failed/absent
# refresh must not wedge the timer. DB-only (no egress), so plain network.target.
After=network.target postgresql.service usa-wa-migrate.service usa-wa-wsl-refresh.service usa-wa-pdc-refresh.service usa-wa-sos-refresh.service
# Alert the operator on a violation (#49/#107). This oneshot IS the anti-drift backstop for a
# missing operator succession event; a silent failure defeats it. %n → this unit's full name.
OnFailure=usa-wa-notify-failure@%n.service

[Service]
Type=oneshot
User=exedev
WorkingDirectory=/home/exedev/usa-wa

# Refuse to start off a non-main checkout (#87).
ExecStartPre=/home/exedev/usa-wa/scripts/assert-main-checkout.sh

# Pin PATH so `uv` resolves regardless of the inherited systemd default.
Environment=PATH=/usr/local/bin:/usr/bin:/bin

# App role only — read-only (SELECT on assignments/roles); no writes, no owner token.
EnvironmentFile=/etc/usa-wa/.env
EnvironmentFile=-/home/exedev/usa-wa/.env

# --frozen --no-sync: run against the installed venv as-is (AGENTS.md § Server Lifecycle, #30).
# Exit 1 on any invariant violation → the OnFailure handler emails the operator.
ExecStart=/usr/local/bin/uv run --frozen --no-sync python -m usa_wa_adapter_legislature.succession_invariants

# A single count query set; bound it so a stuck DB fails deterministically, timer retries next fire.
TimeoutStartSec=300

# This unit is started by usa-wa-succession-invariants.timer; no [Install] WantedBy.
