# deploy/usa-wa-wsl-refresh.service
[Unit]
Description=usa-wa — WSL SOAP refresh (committees + synthesized anchors) one-shot
# usa-wa-migrate.service: schema must be at head before the DML-only refresh
# runs. Ordering only — a failed/absent migrate must not wedge the timer on an
# already-migrated box.
After=network-online.target postgresql.service usa-wa-migrate.service
Wants=network-online.target
# Alert the operator on a failed/timed-out refresh (#49) — on this daily
# unattended timer a failure is otherwise only a journal line nobody watches.
# %n expands to this unit's full name (the notify handler's instance).
OnFailure=usa-wa-notify-failure@%n.service

[Service]
Type=oneshot
User=exedev
WorkingDirectory=/home/exedev/usa-wa

# Pin PATH so `uv` resolves regardless of the inherited systemd default.
Environment=PATH=/usr/local/bin:/usr/bin:/bin

# App role only. refresh.py reads DATABASE_URL (DML — no DDL). DATABASE_URL_OWNER
# also lives in /etc/usa-wa/.env on this single VM but refresh never reads it
# (only alembic does). Later EnvironmentFiles override earlier ones.
EnvironmentFile=/etc/usa-wa/.env
EnvironmentFile=-/home/exedev/usa-wa/.env

# --frozen --no-sync: run against the installed venv as-is; never re-lock or
# sync at unit start. The daily timer therefore can't silently apply a
# dependency change a `git pull` landed in uv.lock — those land only via a
# deliberate `uv sync` in the deploy runbook (AGENTS.md § Server Lifecycle, #30).
ExecStart=/usr/local/bin/uv run --frozen --no-sync python -m usa_wa_adapter_legislature.refresh

# Bound the network-facing zeep WSDL fetch + SOAP call so a hung WSL endpoint
# fails deterministically instead of relying on the manager default; the timer
# retries on the next fire.
TimeoutStartSec=300

# This unit is started by usa-wa-wsl-refresh.timer; no [Install] WantedBy.
