# deploy/usa-wa-sync-powermap.service
[Unit]
Description=usa-wa Power Map sync sidecar — bidirectional canonical-data sync to Power Map
# usa-wa-migrate.service: schema must be at head before this DML-only app-role
# sidecar starts. Redundant with migrate's reciprocal Before=, but declared here
# too for local discoverability. Ordering only — a failed/absent migrate must not
# wedge the sidecar on an already-migrated box (no Wants=/Requires=).
After=network-online.target postgresql.service usa-wa-migrate.service
Wants=network-online.target

[Service]
Type=simple
User=exedev
WorkingDirectory=/home/exedev/usa-wa

# Load system secrets (DATABASE_URL, POWERMAP_API_KEY, POWERMAP_BASE_URL), then
# optional repo-local overrides. Later EnvironmentFiles override earlier ones.
EnvironmentFile=/etc/usa-wa/.env
EnvironmentFile=-/home/exedev/usa-wa/.env

# --frozen --no-sync: run against the installed venv as-is; never re-lock or
# sync at unit start. Dependency changes land only via a deliberate `uv sync`
# in the deploy runbook (AGENTS.md § Server Lifecycle, issue #30).
ExecStart=/usr/local/bin/uv run --frozen --no-sync python -m usa_wa_sync_powermap
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
