# deploy/usa-wa-reconcile-committee-names.timer
[Unit]
Description=usa-wa — weekly committee rename detection
Documentation=https://github.com/CannObserv/usa-wa/issues/53

[Timer]
# Weekly, Sunday 07:30 UTC — staggered 30 min off the active-flag reconcile
# (Sun 07:00) so the two committee-roster oneshots don't contend on the WSL SOAP
# endpoint, and both land after the 06:00 daily WSL refresh has freshened the
# cohort they diff against. Committee renames only occur at biennium boundaries
# (~2 yrs), so weekly is well above the freshness floor; it lands a rollover
# within ≤7 days with no manual action. Reconcile is idempotent — a no-diff run
# emits nothing — so a missed-then-caught run is safe.
OnCalendar=Sun *-*-* 07:30:00 UTC
# Run on next boot if the box was down at the scheduled time.
Persistent=true
# Jitter to avoid firing on the dot.
RandomizedDelaySec=300

[Install]
WantedBy=timers.target
