# deploy/usa-wa-migrate.service
[Unit]
Description=usa-wa — apply Alembic migrations + DB grants under the owner role
After=network.target postgresql.service
# Bring schema to head before the API/sidecar (which run as the DML-only app
# role) start. Ordering only — not a hard dependency; a failed migrate must not
# wedge the boot of an already-migrated box. The edge is declared bilaterally:
# both consumers also carry After=usa-wa-migrate.service, so either side is
# removable without breaking boot ordering (keeps the #50 one-sided-read drift
# from recurring in either direction).
Before=usa-wa.service usa-wa-sync-powermap.service
# Alert the operator if a migration fails (#49). Ordering-only Before= means a
# failed migrate doesn't wedge boot, so the failure is otherwise silent — and a
# half-applied/failed schema is exactly what must not pass unnoticed. %n expands
# to this unit's full name (the notify handler's instance).
OnFailure=usa-wa-notify-failure@%n.service

[Service]
Type=oneshot
RemainAfterExit=yes
User=exedev
WorkingDirectory=/home/exedev/usa-wa

# Refuse to start off a non-main checkout (#87): the enforcement behind the
# "main is the deployed code" convention, which #84 showed is not self-enforcing.
ExecStartPre=/home/exedev/usa-wa/scripts/assert-main-checkout.sh

# Pin PATH so `uv` + `psql` resolve regardless of the inherited systemd default.
Environment=PATH=/usr/local/bin:/usr/bin:/bin

# Owner DSN (DATABASE_URL_OWNER) is read from /etc/usa-wa/.env; alembic/env.py
# prefers it over DATABASE_URL. On this single VM the API + sidecar units load
# the same file, so they also see DATABASE_URL_OWNER in their environment — they
# just never read it (only alembic does). To enforce least privilege on a
# multi-host split, move the owner DSN into a file loaded only by this unit.
EnvironmentFile=/etc/usa-wa/.env
EnvironmentFile=-/home/exedev/usa-wa/.env

ExecStart=/home/exedev/usa-wa/scripts/migrate.sh

[Install]
WantedBy=multi-user.target
