# deploy/usa-wa-sync-powermap.service
[Unit]
Description=usa-wa Power Map sync sidecar — bidirectional canonical-data sync to Power Map
# usa-wa-migrate.service: schema must be at head before this DML-only app-role
# sidecar starts. Redundant with migrate's reciprocal Before=, but declared here
# too for local discoverability. Ordering only — a failed/absent migrate must not
# wedge the sidecar on an already-migrated box (no Wants=/Requires=).
After=network-online.target postgresql.service usa-wa-migrate.service
Wants=network-online.target
# Bound restart attempts so an off-main checkout — which fails the guard
# ExecStartPre below on every attempt — eventually settles into a visible
# `failed` state rather than an unbounded RestartSec=5 loop (#87 CR). Widened
# past systemd's default 10s window (which never trips at RestartSec=5) but kept
# generous so a transient dependency blip still self-heals via Restart=: gives
# up only after 10 failures across 5 min (#87 CR round 2). Invariant asserted in
# test_unit_ordering.py: StartLimitIntervalSec >= RestartSec * StartLimitBurst,
# so the burst can accumulate inside the window and the loop is provably bounded.
StartLimitIntervalSec=300
StartLimitBurst=10

[Service]
Type=simple
User=exedev
WorkingDirectory=/home/exedev/usa-wa

# Pin PATH so the guard's `git` (and `uv`) resolve regardless of the inherited
# systemd default — parity with the oneshot units.
Environment=PATH=/usr/local/bin:/usr/bin:/bin

# Refuse to start off a non-main checkout (#87): the enforcement behind the
# "main is the deployed code" convention, which #84 showed is not self-enforcing.
ExecStartPre=/home/exedev/usa-wa/scripts/assert-main-checkout.sh

# Load system secrets (DATABASE_URL, POWERMAP_API_KEY, POWERMAP_BASE_URL), then
# optional repo-local overrides. Later EnvironmentFiles override earlier ones.
EnvironmentFile=/etc/usa-wa/.env
EnvironmentFile=-/home/exedev/usa-wa/.env

# --frozen --no-sync: run against the installed venv as-is; never re-lock or
# sync at unit start. Dependency changes land only via a deliberate `uv sync`
# in the deploy runbook (AGENTS.md § Server Lifecycle, issue #30).
ExecStart=/usr/local/bin/uv run --frozen --no-sync python -m usa_wa_sync_powermap
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
