# deploy/usa-wa.service
[Unit]
Description=usa-wa — Washington State law, regulation, and policy tracking service
# usa-wa-migrate.service: schema must be at head before this DML-only app-role
# service serves. Redundant with migrate's reciprocal Before=, but declared here
# too for local discoverability. Ordering only — a failed/absent migrate must not
# wedge serving on an already-migrated box (no Wants=/Requires=).
After=network.target postgresql.service usa-wa-migrate.service

[Service]
Type=simple
User=exedev
WorkingDirectory=/home/exedev/usa-wa

# Create runtime dir (runs as root before User= takes effect with + prefix)
ExecStartPre=+/bin/bash -c 'mkdir -p /run/usa-wa && chown exedev:exedev /run/usa-wa'

# Write current git SHA to a runtime env file before starting
ExecStartPre=/bin/bash -c 'echo BUILD_ID=$(git rev-parse --short HEAD) > /run/usa-wa/build-id'

# Load BUILD_ID, then system secrets, then optional repo-local overrides.
# Order matters: later EnvironmentFiles override earlier ones.
EnvironmentFile=-/run/usa-wa/build-id
EnvironmentFile=/etc/usa-wa/.env
EnvironmentFile=-/home/exedev/usa-wa/.env

# --frozen --no-sync: run against the installed venv as-is; never re-lock or
# sync at unit start. Dependency changes land only via a deliberate `uv sync`
# in the deploy runbook (AGENTS.md § Server Lifecycle, issue #30).
ExecStart=/usr/local/bin/uv run --frozen --no-sync uvicorn usa_wa_api.api.main:app --host 0.0.0.0 --port 8000
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
