# deploy/usa-wa-pdc-refresh.service
[Unit]
Description=usa-wa — PDC SODA refresh (House member District + Position) one-shot
# usa-wa-migrate.service: schema must be at head before the DML-only refresh runs.
# usa-wa-wsl-refresh.service: PDC binds Position onto the WSL House Persons, so the
# WSL refresh must have materialized them first. Ordering only — a failed/absent
# predecessor must not wedge the timer on an already-populated box (the normalizer
# logs pdc_house_person_absent and skips an unmatched winner).
After=network-online.target postgresql.service usa-wa-migrate.service usa-wa-wsl-refresh.service
Wants=network-online.target
# Alert the operator on a failed/timed-out refresh (#49).
OnFailure=usa-wa-notify-failure@%n.service

[Service]
Type=oneshot
User=exedev
WorkingDirectory=/home/exedev/usa-wa

# Pin PATH so `uv` resolves regardless of the inherited systemd default.
Environment=PATH=/usr/local/bin:/usr/bin:/bin

# App role only (DML — no DDL). USA_WA_PDC_APP_TOKEN (optional) raises Socrata's
# rate limit; absent is fine. Later EnvironmentFiles override earlier ones.
EnvironmentFile=/etc/usa-wa/.env
EnvironmentFile=-/home/exedev/usa-wa/.env

# --frozen --no-sync: run against the installed venv as-is; never re-lock or sync
# at unit start (AGENTS.md § Server Lifecycle, #30).
ExecStart=/usr/local/bin/uv run --frozen --no-sync python -m usa_wa_adapter_pdc.refresh

# Bound the WSL GetSponsors SOAP call + the PDC SODA GET so a hung endpoint fails
# deterministically; the timer retries on the next fire.
TimeoutStartSec=300

# This unit is started by usa-wa-pdc-refresh.timer; no [Install] WantedBy.
